# backups/systemd/2026-06-18/ollama.service.d/override.conf
[Service]
Environment="OLLAMA_ORIGINS=*"
Environment="OLLAMA_CONTEXT_LENGTH=32768"
Environment="OLLAMA_NUM_PARALLEL=1"
Environment="OLLAMA_MAX_LOADED_MODELS=1"
Environment="OLLAMA_MAX_QUEUE=2"
Environment="OLLAMA_KEEP_ALIVE=5m"
Environment="OLLAMA_FLASH_ATTENTION=1"
Environment="OLLAMA_NOPRUNE=1"
Environment="OLLAMA_NUM_THREADS=20"
MemoryHigh=64G
RestartSec=30
StandardOutput=journal
StandardError=journal
