# com.ura.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ura.dashboard</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/ramonesnaola/URA/ura_ia_1972/.venv/bin/python</string>
        <string>/Users/ramonesnaola/URA/ura_ia_1972/dashboard/ura_web.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/ramonesnaola/URA/ura_ia_1972</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/ramonesnaola/URA/ura_ia_1972/logs/ura_web.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ramonesnaola/URA/ura_ia_1972/logs/ura_web_error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
