# backups/systemd/2026-06-18/llama-vision.service
[Unit]
Description=llama.cpp Vision Model for URA (Qwen2-VL-7B)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ramon
ExecStart=/opt/llama.cpp/build/bin/llama-server --jinja \
    -m /opt/modelos/qwen2-vl-7b-q4_k_m.gguf \
    --mmproj /opt/modelos/qwen2-vl-7b-mmproj-f16.gguf \
    --port 11436 \
    --host 0.0.0.0 \
    --n-gpu-layers 99 \
    -c 8192 \
    --flash-attn 1 \
    --no-jinja \
    --cont-batching
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
