# backups/systemd/2026-06-18/snap-firmware\x2dupdater-227.mount
[Unit]
Description=Mount unit for firmware-updater, revision 227
After=snapd.mounts-pre.target
Before=snapd.mounts.target

[Mount]
What=/var/lib/snapd/snaps/firmware-updater_227.snap
Where=/snap/firmware-updater/227
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=snapd.mounts.target
WantedBy=multi-user.target
