# backups/systemd/2026-06-18/snap-core24-1644.mount
[Unit]
Description=Mount unit for core24, revision 1644
After=snapd.mounts-pre.target
Before=snapd.mounts.target

[Mount]
What=/var/lib/snapd/snaps/core24_1644.snap
Where=/snap/core24/1644
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=snapd.mounts.target
WantedBy=multi-user.target
