# backups/systemd/2026-06-18/snap-mesa\x2d2404-1166.mount
[Unit]
Description=Mount unit for mesa-2404, revision 1166
After=snapd.mounts-pre.target
Before=snapd.mounts.target

[Mount]
What=/var/lib/snapd/snaps/mesa-2404_1166.snap
Where=/snap/mesa-2404/1166
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=snapd.mounts.target
WantedBy=multi-user.target
