# backups/systemd/2026-06-18/ura-network.timer
[Unit]
Description=URA Network Scanner Timer

[Timer]
OnBootSec=30s
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
