# backups/systemd/2026-06-18/ura-model-preliminar.timer
[Unit]
Description=Pre-carga modelo qwen2.5-coder:32b al arranque

[Timer]
OnBootSec=15
Persistent=false

[Install]
WantedBy=timers.target
