# backups/systemd/2026-06-18/ura-maintenance.service
[Unit]
Description=URA Maintenance System - Limpieza automatizada
After=network.target

[Service]
Type=oneshot
User=ramon
Group=ramon
WorkingDirectory=/home/ramon/URA/ura_ia_1972
ExecStart=/usr/bin/python3 /home/ramon/URA/ura_ia_1972/mantenimiento/ura_maintenance.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ura-maintenance

[Install]
WantedBy=multi-user.target
