# backups/systemd/2026-06-18/ura-ejecutor.service
[Unit]
Description=URA Ejecutor API
After=network.target
Wants=network.target

[Service]
Type=simple
User=ramon
WorkingDirectory=/home/ramon/URA/ura_ia_1972
ExecStart=/usr/bin/python3 /home/ramon/URA/ura_ia_1972/scripts/pro/ejecutor_api.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
MemoryHigh=512M
MemoryMax=1G

[Install]
WantedBy=multi-user.target
