# backups/systemd/2026-06-18/ura-auditd-watchdog.service
[Unit]
Description=URA auditd watchdog — check every 5 min
[Service]
Type=oneshot
ExecStart=/usr/local/bin/ura-auditd-watchdog.sh
