# backups/systemd/2026-06-18/ura-hetzner-tunnel.service
[Unit]
Description=URA SSH Tunnel to Hetzner
After=network.target

[Service]
Type=simple
User=ramon
ExecStart=/usr/bin/ssh -o StrictHostKeyChecking=no -o ServerAliveInterval=30 -o ServerAliveCountMax=3 -o ExitOnForwardFailure=yes -N \
    -L 8888:127.0.0.1:8888 \
    -L 3000:127.0.0.1:3000 \
    -D 1080 \
    root@178.105.81.83
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
