# backups/systemd/2026-06-18/ura-detector.service
[Unit]
Description=URA YOLOv8 Detector + ByteTrack + Behavior Analysis
After=network-online.target
Wants=llama-vision.service

[Service]
Type=simple
User=ramon
WorkingDirectory=/opt/ura
Environment=URA_TOKEN=ura_blackwell_2026
Environment=VISION_URL=http://localhost:11436/v1/chat/completions
ExecStart=/usr/bin/python3 /opt/ura/agents/detector_gx10.py
Restart=always
RestartSec=10
StandardOutput=append:/opt/ura/logs/detector.log
StandardError=append:/opt/ura/logs/detector.log

[Install]
WantedBy=multi-user.target
