# backups/systemd/2026-06-18/ura-agent-hierarchy.service
[Unit]
Description=URA Agent Hierarchy System
After=network.target

[Service]
Type=simple
User=ramon
WorkingDirectory=/home/ramon/URA/ura_ia_1972
ExecStart=/usr/bin/python3 /home/ramon/URA/ura_ia_1972/agent_hierarchy.py
Restart=on-failure
RestartSec=5
StartLimitBurst=3
StartLimitIntervalSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
