# backups/systemd/2026-06-18/ura-cleanup.service
[Unit]
Description=URA Cleanup — purge stale files, rotate logs
[Service]
Type=oneshot
ExecStart=/usr/local/bin/ura-cleanup.sh
Nice=19
IOSchedulingClass=idle
