# backups/systemd/2026-06-18/ura-hetzner-watchdog.service
[Unit]
Description=URA Hetzner Watchdog (one-shot)
After=network-online.target tailscaled.service

[Service]
Type=oneshot
ExecStart=/bin/bash /home/ramon/URA/ura_ia_1972/scripts/pro/hetzner_watchdog.sh
User=ramon
EnvironmentFile=/home/ramon/URA/hetzner.env

[Install]
WantedBy=multi-user.target
