# backups/systemd/2026-06-18/ura-mkdocs.service
[Unit]
Description=URA MkDocs — Base de Conocimiento y Autopsias
After=network-online.target
[Service]
Type=simple
User=ramon
WorkingDirectory=/opt/ura
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ramon/.local/bin/mkdocs serve --config-file /opt/ura/mkdocs.yml --dev-addr 0.0.0.0:8088
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
[Install]
WantedBy=multi-user.target
