# backups/systemd/2026-06-18/ura-watchdog.service
[Unit]
Description=URA Watchdog Service

[Service]
Type=oneshot
ExecStart=/home/ramon/ura-fixes/ura-watchdog.sh
User=ramon
