# backups/systemd/2026-06-18/ura-ufw-rules.service
[Unit]
Description=URA UFW Tailscale rules (una vez al boot)
After=network.target ufw.service
Before=tailscaled.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ura-ufw-rules.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
