# backups/systemd/2026-06-18/ura-watchdog-buffer.service
[Unit]
Description=URA Watchdog de Buffer de 30GB
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash /home/ramon/URA/ura_ia_1972/scripts/pro/watchdog_buffer.sh
Restart=always
RestartSec=10
User=ramon

[Install]
WantedBy=multi-user.target
