# backups/systemd/2026-06-18/ura-pipeline.service
[Unit]
Description=URA Motor de Conocimiento - Pipeline GX10
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ura-motor --config /etc/ura/config.json pipeline
Environment=URA_CONFIG=/etc/ura/config.json
User=root
StandardOutput=journal
StandardError=journal
