# backups/systemd/2026-06-18/ura-memory-watchdog.service
[Unit]
Description=URA Memory Pressure Watchdog

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ura-memory-watchdog.sh
User=ramon
