# backups/systemd/2026-06-18/ura-ssh-guard.service
[Unit]
Description=URA SSH Guard
After=network-online.target tailscaled.service

[Service]
Type=simple
ExecStart=/usr/local/bin/ura-ssh-guard.sh
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
