# backups/systemd/2026-06-18/ura-procesamiento-lento.service
[Unit]
Description=URA Daemon de Procesamiento Lento (10% CPU)
After=network.target
Wants=ura-watchdog-buffer.service

[Service]
Type=simple
ExecStart=/bin/bash /home/ramon/URA/ura_ia_1972/scripts/pro/daemon_procesamiento_lento.sh
Restart=always
RestartSec=30
User=ramon
CPUSchedulingPolicy=idle
IOSchedulingClass=idle

[Install]
WantedBy=multi-user.target
