# backups/systemd/2026-06-18/ura-mochila.service.d/log_policy.conf
[Service]
StandardOutput=journal
StandardError=journal
LogRateLimitIntervalSec=30
LogRateLimitBurst=100
