# scripts/pro/com.ura.voice.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://apple.com">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ura.voice</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/ramonesnaola/URA/ura_ia_1972/scripts/pro/demo_pipeline_mac.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ramonesnaola/URA/ura_ia_1972</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/ramonesnaola/URA/ura_ia_1972/logs/mac_voice_out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ramonesnaola/URA/ura_ia_1972/logs/mac_voice_err.log</string>
</dict>
</plist>
