# scripts/pro/gx10-api.service
[Unit]
Description=URA GX10 API — Remote endpoint with post-crash audit gate
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
User=ramon
WorkingDirectory=/opt/ura
EnvironmentFile=/etc/ura/gx10-api.env
Environment=PYTHONUNBUFFERED=1

# Ejecutar auditoría post-caída antes de arrancar la API.
# Si la auditoría detecta corrupción, el servicio falla y se reporta al panel.
ExecStartPre=/opt/ura/scripts/auditoria_post_caida.sh
ExecStart=/usr/bin/python3 /opt/ura/agents/ura_api.py

Restart=on-failure
RestartSec=30
StartLimitBurst=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
