# backups/systemd/2026-06-18/openclaw.service.d/override.conf
[Service]
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=5
