# scripts/legacy/start_whisper_user.service
[Unit]
Description=URA Whisper Speech-to-Text Server
After=default.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/ramon/URA/ura_ia_1972/services/whisper_server.py
WorkingDirectory=/home/ramon/URA/ura_ia_1972
Environment=PYTHONPATH=/home/ramon/URA/ura_ia_1972
Restart=on-failure
RestartSec=5
StandardOutput=append:/home/ramon/URA/ura_ia_1972/logs/whisper_server.log
StandardError=append:/home/ramon/URA/ura_ia_1972/logs/whisper_server.log

[Install]
WantedBy=default.target
