# scripts/legacy/start_router_user.service
[Unit]
Description=URA Router llama_router (usuario)
After=default.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/ramon/URA/ura_ia_1972/services/llama_router.py --host 100.127.206.86 --models codestral-22b qwen2.5-coder-q8 qwen2.5-coder-32b
WorkingDirectory=/home/ramon/URA/ura_ia_1972
Environment=PYTHONPATH=/home/ramon/URA/ura_ia_1972
Environment=OLLAMA_MAX_LOADED_MODELS=4
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/ramon/URA/ura_ia_1972/logs/llama_router.log
StandardError=append:/home/ramon/URA/ura_ia_1972/logs/llama_router.log

[Install]
WantedBy=default.target
