# backups/systemd/2026-06-18/tailscaled.service.d/override.conf
[Unit]
Description=Tailscale node agent (URA hardened)
After=network-online.target
Wants=network-online.target

[Service]
Restart=always
RestartSec=5s
