# deploy/model-router.service
[Unit]
Description=URA Model Router (intelligent model selection)
After=network.target ollama.service

[Service]
Type=simple
User=ramon
WorkingDirectory=/home/ramon/URA/ura_ia_1972
Environment="URA_AUTH_ENABLED=false"
ExecStartPre=/bin/bash -c '/bin/fuser -k 11435/tcp 2>/dev/null || true'
ExecStart=/usr/bin/python3 -m core.model_router
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
