# deploy/opencode.service
[Unit]
Description=OpenCode Web Server with OpenClaw Gateway and Ollama
After=network.target openclaw-gateway.service

[Service]
Type=simple
User=ramon
WorkingDirectory=/home/ramon/URA/ura_ia_1972
Environment="OPENCODE_GATEWAY_URL=ws://10.164.1.99:18789"
Environment="OPENCODE_GATEWAY_TOKEN=V0T9kzrReQyvUlcPGrzlqO-BNVM_xIc1oJN0KZhftrU"
Environment="OPENCODE_SERVER_USERNAME=ramon"
Environment="OPENCODE_SERVER_PASSWORD=ura_1972_secure_autonomous"
ExecStart=/home/ramon/.opencode/bin/opencode web --port 8081 --hostname 0.0.0.0
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
