# deploy/ura-router-health.service
[Unit]
Description=URA Model Router Health Check
Documentation=https://github.com/anomalyco/ura

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/ura-router-health.sh
StandardOutput=journal
StandardError=journal
