# deploy/ura-router-health.service
[Unit]
Description=URA Model Router Health Check
Documentation=https://github.com/anomalyco/ura

[Service]
Type=oneshot
User=root
ExecStart=/home/ramon/URA/ura_ia_1972/scripts/pro/health_check_router.sh
StandardOutput=journal
StandardError=journal
