# deploy/ura-router-alerts.service
[Unit]
Description=URA Model Router Alert Monitor
Documentation=https://github.com/anomalyco/ura

[Service]
Type=oneshot
User=root
ExecStart=/bin/bash -c 'journalctl -u model-router.service --since "1 hour ago" | grep -c "Failed with result" | awk "{if(\$1>10) exit 1}"'
StandardOutput=journal
StandardError=journal
