# deploy/ura-search-pipeline.timer
[Unit]
Description=URA-Search: Pipeline diario (cada 12h)
Documentation=https://github.com/anomalyco/ura

[Timer]
OnCalendar=*-*-* 00,12:00:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
