# deploy/ura-search-pipeline.service
[Unit]
Description=URA-Search: Pipeline completo de recolección
Documentation=https://github.com/anomalyco/ura

[Service]
Type=oneshot
User=ramon
WorkingDirectory=/home/ramon/URA/ura_ia_1972
ExecStart=/bin/bash -c '
  echo "=== URA-Search Pipeline ==="
  date
  echo "Fase 1: Recolección (collector_base)"
  python3 -m ura_search.collector_base batch --batch ura_search/config/fuentes.json --workers 5
  echo "Fase 2: Agentes Docker"
  bash agents/launch-agents.sh
  echo "Fase 3: Enriquecimiento IA"
  python3 -m ura_search.enricher_text --all
  echo "Fase 4: Indexación"
  python3 -m ura_search.indexer --index
  echo "=== Pipeline completado ==="
'
StandardOutput=append:/home/ramon/URA/ura_ia_1972/logs/pipeline.log
StandardError=append:/home/ramon/URA/ura_ia_1972/logs/pipeline.log
