# backups/systemd/2026-06-18/nvidia-persistenced.service.d/nv-persistence-override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/nvidia-persistenced --persistence-mode --verbose

[Install]
WantedBy=
WantedBy=basic.target
