# backups/systemd/2026-06-18/ura-executor.service.d/execstart.conf
[Service]
ExecStart=
ExecStart=/usr/bin/python3 /home/ramon/URA/ura_ia_1972/scripts/pro/ejecutor_api.py
WorkingDirectory=/home/ramon/URA/ura_ia_1972
Environment=OLLAMA_URL=http://localhost:11434
Environment=EXECUTOR_HOST=0.0.0.0
