# backups/systemd/2026-06-18/ura-executor.service.d/recursos.conf
[Service]
MemoryMax=1G
CPUWeight=50
MemoryHigh=768M
IOWeight=25
