# backups/systemd/2026-06-18/ura-executor.service.d/workdir.conf
[Service]
WorkingDirectory=/home/ramon/URA/ura_ia_1972
