# backups/systemd/2026-06-18/docker.service.d/override.conf
[Unit]
After=nvidia-gpu-reset.target
Wants=nvidia-gpu-reset.target
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd -H fd://
